# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn --bind 0.0.0.0:9191 app:app

[program:all]
command=/bin/bash /start-spark.sh